# home/.bash/profile.d/dotsplat.sh
source "$HOME/.dotsplat.git/completions/dotsplat-completion.bash"
source "$HOME/.dotsplat.git/dotsplat.sh"
dotsplat --quiet refresh
