# home/.bash/profile.d/prompt.sh
source $HOME/.dotsplat/repos/dotfiles/home/.bash/prompt/git-svn-bash.sh

